# cookbooks/ey-bin/templates/default/activemq.erb
#!/bin/bash
#
# chkconfig: 345 99 01
# description: Message queue
# This script should live in /engineyard/bin/
# You should call this script from monit
PROG=activemq
DAEMON_HOME=/data/$1/activemq
DAEMON=$DAEMON_HOME/bin/$PROG
LOCKFILE=/var/run/$PROG/$1.lock
PIDFILE=/var/run/$PROG/$1.pid

test -x $DAEMON || exit 0

# Source function library.
. /etc/init.d/functions.sh

RETVAL=0

usage () {
    echo "Usage: /engineyard/bin/activemq zaadz {start|stop|restart|status}"
    RETVAL=1
}

start () {
  echo -n $"Starting $PROG: "
  PID=`cat $PIDFILE`
  if [ ! -e /proc/${PID} ] ; then
    rm -f $LOCKFILE
  fi
  if [ ! -e $LOCKFILE ]; then
    cd $DAEMON_HOME
    $DAEMON > >(logger -t $PROG) 2>&1 &
  else
    echo -n "Lockfile exists"
    false
  fi
  RETVAL=$?
  if [ $RETVAL -eq 0 ]; then
    logger -t activemq "starting $PROG."
    echo $! > $PIDFILE
    touch $LOCKFILE
  else
    logger -t activemq "unable to start $PROG."
  fi
  echo
}

stop () {
  echo "Shutting down $PROG/$1: "
  PID=`cat $PIDFILE`
  if [ -e /proc/${PID} ] ; then
    kill $PID
    sleep 5
    if [ -e /proc/${PID} ] ; then
      kill -9 $PID
      sleep 5
      if [ -e /proc/${PID} ] ; then
        echo "Process ${PID} won't die!"
        exit 1
      fi
    fi
  fi
  rm -f $LOCKFILE
}

case "$2" in
  start) start ;;
  stop) stop ;;
  restart|reload)
    stop
    start
    ;;
  status)
    status $PROG -p $PIDFILE
    RETVAL=$?
    ;;
  *) usage ;;
esac

exit $RETVAL
